# templates/etc/systemd/system/docker.service.d/environment.conf.j2
# {{ ansible_managed }}

[Service]
{% for env in docker__daemon_environment %}
Environment="{{ env }}"
{% endfor %}
