# templates/etc/systemd/system/docker.service.d/environment.conf.j2
# {{ ansible_managed }}

[Service]
Environment="{{ docker__daemon_environment | join('" "') }}"
